# templates/nimbus-eth1-long-benchmark.service.j2
[Unit]
Description=Nimbus Eth1 short import benchmarking service (1 week)
Documentation=https://github.com/status-im/nimbus-eth1
Requires=network-online.target
After=network-online.target

[Service]
User={{ nimbus_eth1_service_user }}
Group={{ nimbus_eth1_service_group }}
WorkingDirectory={{ nimbus_eth1_service_path }}
SyslogIdentifier=nimbus-eth1-{{ nimbus_eth1_network }}-short-benchmark
SyslogFacility=local6
SyslogLevel=debug
RuntimeMaxSec=168hr
ExecStart={{ nimbus_eth1_repo_path }}/build/nimbus_execution_client import \
    --network={{ nimbus_eth1_network }} \
    --data-dir='{{ nimbus_eth1_data_path }}/{{ nimbus_eth1_data_folder }}' \
    --era1-dir={{ nimbus_eth1_era1_dir | mandatory }} \
    --era-dir={{ nimbus_eth1_era_dir | mandatory }}

[Install]
WantedBy=multi-user.target
